# docker/var/spool/crontabs/nginx
0 * * * * cd /var/www/app && php cronjob.php
